# packages/graphics/bcm2835-driver/system.d/fbset.service
[Unit]
Description=Framebuffer setup
After=kodi.service

[Service]
Type=oneshot
ExecStart=/usr/sbin/fbset -g 1 1 1 1 32
RemainAfterExit=yes

[Install]
WantedBy=kodi.target
