# packages/tools/plymouth/system.d/plymouth-quit.service
[Unit]
Description=Terminate Plymouth Boot Screen
After=plymouth-start.service

[Service]
ExecStart=-/usr/bin/plymouth quit --retain-splash
Type=oneshot
TimeoutSec=20
